# playbook-redmine/roles/redmine/templates/redmine-nginx.conf.j2
[program:redmine-nginx]
command={{ nginx_dir }}/sbin/nginx
autostart=true
autorestart=true
stopsignal=QUIT

